# tun-bridge-client.service
[Unit]
Description=TUN Network Bridge Client Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/tun_bridge --client --host SERVER_IP --port 8080 --ip 10.0.1.2 --remote-ip 10.0.1.1 --routes "192.168.2.0/24,172.16.0.0/12"
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tun-bridge-client

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/dev/net/tun /proc/sys/net
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

[Install]
WantedBy=multi-user.target
